# distrib/etc/sv/dbus/run
#!/bin/sh

if [ -f ./svc.conf ]; then
	. ./svc.conf
fi

exec 2>&1

chown messagebus:messagebus /var/run/dbus

exec chpst -u messagebus /usr/local/bin/dbus-daemon --nofork $flags
