# distrib/etc/sv/sshd/run
#!/bin/sh

if [ -f ./svc.conf ]; then
	. ./svc.conf
fi

exec 2>&1
# the -d option prevents forking
exec /usr/sbin/sshd -D -e $flags

