# distrib/etc/sv/moused/run
#!/bin/sh

moused -f -p /dev/psm0 -t auto
