# proto/etc/sv/amd/run
#!/bin/sh

/sbin/sv start rpcbind >/dev/null || exit 1

sleep 0.5 # give rpcbind a chance to start properly

. ./svc.conf

exec 2>&1
exec 1>/dev/null

exec /usr/sbin/amd -N $flags
