# proto/etc/sv/sshd/finish
#!/bin/sh

/bin/rm -f /var/run/sshd.pid
