# distrib/etc/sv/rpcbind/finish
#!/bin/sh

/bin/rm -f /var/run/rpcbind.*
